# common/skills/spinup/scripts/inc.example.cmux-close-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>inc.example.cmux-close-listener</string>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>ProgramArguments</key>
	<array>
		<string>$HOME/.claude/skills/spinup/scripts/cmux-close-listener.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>$HOME/.claude/cache/cmux-close-listener.log</string>
	<key>StandardOutPath</key>
	<string>$HOME/.claude/cache/cmux-close-listener.log</string>
</dict>
</plist>
